# misc/dropbatch.path
[Unit]
Description=Dropbatch job monitor

[Path]
PathExistsGlob=%h/Dropbox/%H-*.sh
PathExistsGlob=%h/Dropbox/.System/Batch/%H-*.sh

[Install]
WantedBy=paths.target
