# misc/dropbatch.timer
[Unit]
Description=Dropbatch job timer

[Timer]
OnCalendar=*:00/15:00

[Install]
WantedBy=timers.target
