# system/arch/kernel-post-upgrade@.path
[Unit]
Description=Kernel post-upgrade watch (for linux-%i)

[Path]
PathChanged=/boot/vmlinuz-linux-%i
PathChanged=/boot/intel-ucode.img
PathChanged=/boot/initramfs-linux-%i.img
PathChanged=/etc/kernel/cmdline
PathChanged=/etc/os-release

[Install]
WantedBy=paths.target
