# system/arch/kernel-post-upgrade.service
[Unit]
Description=Kernel post-upgrade script

[Service]
Type=oneshot
ExecStart=/boot/kernel-post-upgrade.sh

# vim: ft=desktop
