# system/arch/kernel-post-upgrade@.service
[Unit]
Description=Kernel post-upgrade script (for %i)

[Service]
Type=oneshot
ExecStart=/boot/kernel-post-upgrade.sh %i
