# system/arch/kernel-post-upgrade@.path
[Unit]
Description=Kernel post-upgrade watch (for %i)

[Path]
PathChanged=/boot/vmlinuz-%i
# TODO: locking to avoid scripts running over each other whenever ucode is updated
#PathChanged=/boot/intel-ucode.img
PathChanged=/boot/initramfs-%i.img
PathChanged=/etc/kernel/cmdline

[Install]
WantedBy=multi-user.target
DefaultInstance=linux
