# system/arch/mkloaderconf.path
# vim: ft=systemd
[Unit]
Description=Boot Kernel Configuration Watch

[Path]
PathChanged=/boot/vmlinuz-*
PathChanged=/boot/initramfs-*
PathChanged=/boot/intel-ucode.img
PathChanged=/etc/kernel/cmdline

[Install]
WantedBy=multi-user.target
