# system/arch/mkloaderconf.path
# vim: ft=systemd
[Unit]
Description=Boot Kernel Configuration Watch

[Path]
PathChanged=/boot
PathChanged=/etc/kernel/cmdline

[Install]
WantedBy=multi-user.target
