# system/arch/mkloaderconf.service
# vim: ft=systemd
[Unit]
Description=Boot Kernel Configuration Update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mkloaderconf
