# kerberos/krb5-creds@.service
[Unit]
Description=Kerberos machine credentials (%i instance)
After=krb5-kdc.service
Wants=remote-fs-pre.target
Before=remote-fs-pre.target

[Service]
Type=forking
ExecStart=/home/grawity/code/kerberos/host-kinit -s %i

[Install]
WantedBy=multi-user.target
